# deploy/autovideosrt.service
[Unit]
Description=AutoVideoSrt Web Service
After=network.target

[Service]
User=root
WorkingDirectory=/opt/autovideosrt
Environment="PATH=/opt/autovideosrt/venv/bin:/usr/bin:/usr/local/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="AUTOVIDEOSRT_GUNICORN_BIND=0.0.0.0:80"
Environment="AUTOVIDEOSRT_GUNICORN_THREADS=32"
Environment="OMP_NUM_THREADS=1"
Environment="OPENBLAS_NUM_THREADS=1"
Environment="MKL_NUM_THREADS=1"
Environment="NUMEXPR_NUM_THREADS=1"
ExecStart=/opt/autovideosrt/venv/bin/gunicorn --config /opt/autovideosrt/deploy/gunicorn.conf.py main:app
Restart=always
RestartSec=5
# Cooperative shutdown is wired in deploy/gunicorn.conf.py:
# post_worker_init chains SIGTERM into the process shutdown coordinator
# and stops APScheduler before Gunicorn drains. Tracked tasks honour the
# cancellation points; worker_exit waits up to 10s for them to drain.
# 45s here = 15s gunicorn graceful_timeout + systemd buffer; publishing must
# still run the active task pre-restart check before restarting the service.
TimeoutStopSec=45
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
