# deploy/tabcut-daily-selection.service
[Unit]
Description=AutoVideoSrt Tabcut daily US selection crawler
After=network-online.target graphical-session.target
Wants=network-online.target

[Service]
Type=oneshot
User=cjh
Group=cjh
WorkingDirectory=/opt/autovideosrt
Environment=PYTHONUNBUFFERED=1
Environment=TABCUT_CDP_URL=http://127.0.0.1:9227
Environment=TABCUT_REGION=US
ExecStart=/opt/autovideosrt/venv/bin/python -m tools.tabcut_crawler.main --days 7

[Install]
WantedBy=multi-user.target
