# deploy/server_browser/autovideosrt-dxm03-rjc-vnc.service
[Unit]
Description=AutoVideoSrt DXM03-RJC visible browser + noVNC runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cjh
Group=cjh
WorkingDirectory=/opt/autovideosrt
Environment=PATH=/opt/autovideosrt/venv/bin:/usr/bin:/usr/local/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=DXM_NAME=DXM03-RJC
Environment=DXM_DISPLAY_NUM=25
Environment=DXM_PROFILE_DIR=/data/autovideosrt/browser/profiles/rjc-dianxiaomi
Environment=DXM_RUNTIME_DIR=/data/autovideosrt/browser/runtime-rjc-dianxiaomi
Environment=DXM_LOG_DIR=/data/autovideosrt/browser/logs/rjc-dianxiaomi
Environment=DXM_CDP_PORT=9225
Environment=DXM_VNC_PORT=5925
Environment=DXM_NOVNC_PORT=6095
Environment=DXM_START_URL=https://www.dianxiaomi.com/web/order/all
ExecStart=/opt/autovideosrt/deploy/server_browser/run_visible_dxm_env.sh
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
