# deploy/server_browser/autovideosrt-shopifyid-sync.service
[Unit]
Description=AutoVideoSrt Shopify ID Dianxiaomi sync
Wants=autovideosrt-dxm03-rjc-vnc.service
After=network-online.target autovideosrt-dxm03-rjc-vnc.service mysql.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/autovideosrt
Environment=SHOPIFYID_BROWSER_MODE=server-cdp
Environment=SHOPIFYID_BROWSER_CDP_URL=http://127.0.0.1:9225
Environment=SHOPIFYID_DB_MODE=local
ExecStart=/opt/autovideosrt/venv/bin/python /opt/autovideosrt/tools/shopifyid_dianxiaomi_sync.py --skip-login-prompt --browser-mode server-cdp --browser-cdp-url http://127.0.0.1:9225 --db-mode local
TimeoutStartSec=7200
