# deploy/server_browser/autovideosrt-roas-backfill.timer
[Unit]
Description=Run AutoVideoSrt ROAS fields backfill every day

[Timer]
OnCalendar=*-*-* 13:07:00
Persistent=true
Unit=autovideosrt-roas-backfill.service

[Install]
WantedBy=timers.target
