# deploy/server_browser/autovideosrt-mingkong-fine-ai-worker.service
[Unit]
Description=AutoVideoSrt Mingkong Fine AI evaluation worker pool
After=network-online.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/autovideosrt
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/autovideosrt/venv/bin/python tools/mingkong_fine_ai_auto_evaluation_worker.py --workers 2 --idle-sleep-seconds 10
Restart=always
RestartSec=10
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
