# deploy/server_browser/autovideosrt-sku-actual-roas.timer
[Unit]
Description=Run AutoVideoSrt SKU actual breakeven ROAS snapshot daily

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true
Unit=autovideosrt-sku-actual-roas.service

[Install]
WantedBy=timers.target
