# deploy/server_browser/autovideosrt-roi-realtime-sync.service
[Unit]
Description=AutoVideoSrt real ROAS realtime day snapshot sync
After=network-online.target autovideosrt-dxm01-meta-vnc.service autovideosrt-dxm03-rjc-vnc.service
Wants=network-online.target autovideosrt-dxm01-meta-vnc.service autovideosrt-dxm03-rjc-vnc.service

[Service]
Type=oneshot
WorkingDirectory=/opt/autovideosrt
EnvironmentFile=-/opt/autovideosrt/.env
Environment=PYTHONUNBUFFERED=1
Environment=META_REALTIME_SYNC_CHANNEL=browser
Environment=META_AD_EXPORT_CDP_URL=http://127.0.0.1:9222
Environment=DXM_ORDER_BROWSER_CDP_URL=http://127.0.0.1:9225
ExecStart=/opt/autovideosrt/venv/bin/python /opt/autovideosrt/tools/roi_hourly_sync.py --lookback-hours 3 --max-scan-pages 40 --meta-channel browser
User=root
Group=root
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800
