# deploy/server_browser/autovideosrt-dxm02-mk-vnc.service
[Unit]
Description=AutoVideoSrt DXM02-MK visible browser + noVNC runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cjh
Group=cjh
WorkingDirectory=/opt/autovideosrt
Environment=PATH=/opt/autovideosrt/venv/bin:/usr/bin:/usr/local/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=DXM_NAME=DXM02-MK
Environment=DXM_DISPLAY_NUM=23
Environment=DXM_PROFILE_DIR=/data/autovideosrt/browser/profiles/mk-selection
Environment=DXM_RUNTIME_DIR=/data/autovideosrt/browser/runtime-mk-selection
Environment=DXM_LOG_DIR=/data/autovideosrt/browser/logs/mk-selection
Environment=DXM_CDP_PORT=9223
Environment=DXM_VNC_PORT=5923
Environment=DXM_NOVNC_PORT=6093
Environment=DXM_START_URL=https://www.dianxiaomi.com/web/stat/salesStatistics
ExecStart=/opt/autovideosrt/deploy/server_browser/run_visible_dxm_env.sh
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
