# deploy/server_browser/autovideosrt-roi-realtime-sync.timer
[Unit]
Description=Run AutoVideoSrt real ROAS realtime day snapshot sync every 20 minutes

[Timer]
OnCalendar=*:02/20
Persistent=true
Unit=autovideosrt-roi-realtime-sync.service

[Install]
WantedBy=timers.target
