# deploy/server_browser/autovideosrt-meta-daily-final-check.timer
[Unit]
Description=Confirm AutoVideoSrt Meta final daily ad data sync at 19:00 Beijing time

[Timer]
OnCalendar=*-*-* 19:00:00
Persistent=true
Unit=autovideosrt-meta-daily-final-check.service

[Install]
WantedBy=timers.target
