# deploy/server_browser/autovideosrt-sku-actual-roas.service
[Unit]
Description=AutoVideoSrt SKU actual breakeven ROAS daily snapshot
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/autovideosrt
EnvironmentFile=-/opt/autovideosrt/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/autovideosrt/venv/bin/python /opt/autovideosrt/tools/sku_actual_roas_snapshot.py
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800
