# deploy/server_browser/autovideosrt-roi-realtime-sync.timer
[Unit]
Description=Run AutoVideoSrt real ROAS realtime day snapshot sync every 20 minutes (aligned with order-profit-incremental)

[Timer]
OnCalendar=*-*-* *:00,20,40:00
Persistent=true
Unit=autovideosrt-roi-realtime-sync.service

[Install]
WantedBy=timers.target
