# deploy/server_browser/autovideosrt-dianxiaomi-listing-ranking-sync.timer
[Unit]
Description=Run AutoVideoSrt Dianxiaomi full recent-sales Listing archive sync every day

[Timer]
OnCalendar=*-*-* 12:40:00
Persistent=true
Unit=autovideosrt-dianxiaomi-listing-ranking-sync.service

[Install]
WantedBy=timers.target
