# deploy/server_browser/autovideosrt-mingkong-product-library-sync.timer
[Unit]
Description=Run AutoVideoSrt Mingkong product library sync weekly

[Timer]
OnCalendar=Mon *-*-* 03:30:00
Persistent=true
Unit=autovideosrt-mingkong-product-library-sync.service

[Install]
WantedBy=timers.target
