# deploy/server_browser/autovideosrt-mingkong-sku-backfill-plan.timer
[Unit]
Description=Run AutoVideoSrt Mingkong weekly SKU backfill planning

[Timer]
OnCalendar=Mon *-*-* 07:30:00
Persistent=true
Unit=autovideosrt-mingkong-sku-backfill-plan.service

[Install]
WantedBy=timers.target
