# deploy/server_browser/autovideosrt-mk-browser.service
[Unit]
Description=AutoVideoSrt MK Selection Isolated Browser Runtime
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/autovideosrt
Environment="PATH=/opt/autovideosrt/venv/bin:/usr/bin:/usr/local/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="APP_DIR=/opt/autovideosrt"
Environment="VENV_DIR=/opt/autovideosrt/venv"
Environment="PLAYWRIGHT_BROWSERS_PATH=/opt/autovideosrt/.playwright-browsers"
EnvironmentFile=-/etc/default/autovideosrt-mk-browser
ExecStart=/opt/autovideosrt/deploy/server_browser/run_server_browser.sh
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
