# deploy/server_browser/autovideosrt-roi-realtime-sync.timer
[Unit]
Description=Run AutoVideoSrt real ROAS realtime day snapshot sync every hour

[Timer]
OnCalendar=*:02:00
Persistent=true
Unit=autovideosrt-roi-realtime-sync.service

[Install]
WantedBy=timers.target
