# deploy/server_browser/autovideosrt-order-profit-incremental.timer
[Unit]
Description=Run AutoVideoSrt order-profit incremental sync every 20 minutes (aligned with roi-realtime-sync)

[Timer]
OnCalendar=*-*-* *:00,20,40:00
Persistent=true
Unit=autovideosrt-order-profit-incremental.service

[Install]
WantedBy=timers.target
