# deploy/server_browser/autovideosrt-mingkong-material-daily-snapshot.service
[Unit]
Description=AutoVideoSrt Mingkong daily material snapshot and yesterday Top100
After=network-online.target mysql.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/autovideosrt
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/autovideosrt/venv/bin/python tools/mingkong_material_daily_snapshot.py --sleep-after-products 10 --sleep-seconds 10
TimeoutStartSec=21600
