# deploy/server_browser/autovideosrt-roi-realtime-sync.service.d/10-browser-lock.conf
# Docs-anchor: docs/superpowers/specs/2026-05-09-roi-hourly-sync-lock-recovery.md
#
# Wrap roi_hourly_sync with the shared server-browser lock so that long
# running browser tasks (roas-backfill, sku-aggregates, etc.) can not collide
# with the hourly ROI tick. Timeout 300s = fail fast: a single tick gives up
# after 5 min so the next tick (every hour) can still pick up.
[Service]
Environment=BROWSER_AUTOMATION_LOCK_TIMEOUT_SECONDS=300
Environment=BROWSER_AUTOMATION_LOCK_RETRY_SECONDS=10
Environment=BROWSER_AUTOMATION_LOCK_ALERT_TASK_CODE=roi_hourly_sync
ExecStart=
ExecStart=/opt/autovideosrt/deploy/server_browser/with_browser_lock.sh /opt/autovideosrt/venv/bin/python /opt/autovideosrt/tools/roi_hourly_sync.py --lookback-hours 3 --max-scan-pages 40 --meta-channel browser
